# supervisor.ini
[program:flask]
command = gunicorn run:app
autostart = true
autorestart = true
startretries = 3


[program:huey]
command = huey_consumer app.huey.run -w 32
autostart = true
autorestart = true
startretries = 3